# systemd/epp-server.service
[Unit]
Description=EPP Server for .AE Domain Registry
Documentation=https://www.rfc-editor.org/rfc/rfc5730
After=network.target oracle-instantclient.service
Wants=network-online.target

[Service]
Type=simple
User=epp
Group=epp

# Working directory
WorkingDirectory=/opt/epp-server

# Use bundled Python virtual environment
Environment="PATH=/opt/epp-server/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/epp-server"

# Oracle Instant Client (adjust path as needed)
Environment="LD_LIBRARY_PATH=/usr/lib/oracle/21/client64/lib"
Environment="ORACLE_HOME=/usr/lib/oracle/21/client64"

# Start command
ExecStart=/opt/epp-server/venv/bin/python -m src.server --config /opt/epp-server/config/epp.yaml

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Allow binding to port 700
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Read-write paths
ReadWritePaths=/opt/epp-server/logs
ReadWritePaths=/opt/epp-server/run

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=epp-server

# Limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
